# TheFrame/sidecar/theframe.service
[Unit]
Description=TheFrame Samsung Frame TV Sidecar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/theframe
ExecStart=/opt/theframe/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Let the TV settle on the network before connecting
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target
